# supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/app/logs
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# Example: Filesystem MCP Server
[program:mcp-filesystem]
command=npx -y @modelcontextprotocol/server-filesystem /tmp
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/logs/mcp-filesystem.log
stderr_logfile=/app/logs/mcp-filesystem-error.log

# Example: GitHub MCP Server
[program:mcp-github]
command=npx -y @modelcontextprotocol/server-github
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/logs/mcp-github.log
stderr_logfile=/app/logs/mcp-github-error.log
; Tokens should be provided via the container environment (.env).
